# systemd/wallpaper.service
[Service]
Type=oneshot
ExecStart=feh -z --bg-scale --no-fehbg /home/sources/wallpaper

